# docker/runit/etc/service/ikea-monitor-fetch/run
#!/bin/bash
cd /srv/ikea-monitor
exec 2>&1
exec chpst -U www-user node src/fetch.js
